# deploy/poly-arb-bot.service
[Unit]
Description=Poly Arb Bot C++ shadow simulation loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/poly-arb-bot
EnvironmentFile=/opt/poly-arb-bot/.env
ExecStartPre=/bin/bash /opt/poly-arb-bot/scripts/check_ntp.sh
ExecStart=/bin/bash /opt/poly-arb-bot/scripts/run_shadow_loop.sh
Restart=always
RestartSec=5
StandardOutput=append:/var/log/poly-arb-bot.log
StandardError=append:/var/log/poly-arb-bot.err.log

[Install]
WantedBy=multi-user.target
